# ops/systemd/sis-db-backup.service
[Unit]
Description=SIS Database Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{USER}}
Group={{GROUP}}
WorkingDirectory={{RUNTIME_ROOT}}
EnvironmentFile={{ENV_FILE}}
Environment=DOTENV_CONFIG_PATH={{ENV_FILE}}
ExecStart={{NODE_BIN}} -r dotenv/config {{RUNTIME_ROOT}}/tools/db-backup-failsafe.mjs
